# apps/daily-log/daily-log.service
[Unit]
Description=Daily Log App
After=network.target

[Service]
Type=simple
User=dima
WorkingDirectory=/home/dima/daily-log
EnvironmentFile=/home/dima/daily-log/.env
ExecStart=/usr/bin/node /home/dima/daily-log/dist/server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
